# package/usbip/usbip.rules
ACTION=="add",SUBSYSTEM=="usb"    RUN+="/usr/sbin/usbip bind -b $kernel"
ACTION=="remove",SUBSYSTEM=="usb" RUN+="/usr/sbin/usbip unbind -b $kernel"
